# mattermost.conf
[program:mattermost]
user = nobody
command = /scripts/start.sh
