# mattermost.conf
[program:mattermost]
user = nobody
command = /opt/mattermost/bin/platform -config="/mattermost/config.json"
